# crontab
#Run every hour
0 * * * * sudo bash ~/PUPget.sh
