# crontab
#Run every 24 hours
0 0 * * * sudo bash ~/PUPgetBackend.sh
