# crontab
#Run every hour
0 * * * * sudo bash ~/PUPgetBackend.sh
